# _development/vagrant/cronjobs/crontab
*/10 * * * * /home/vagrant/cronjobs/logrotate.sh
